# RaspberryPi/systemd/odrive-can-suboptimal.service
[Unit]
Description=ODrive can_runner service (starts the controller and tears down can0 on stop)
After=multi-user.target network.target dev-gpiomem.device dev-spidev0.0.device
Requires=network.target

[Service]
Type=simple
User=neonaero
WorkingDirectory=/home/neonaero/Documents/odrive-controller/odrive-controller-odrive-rpi/RaspberryPi
# Bring up can0 only while this service runs. Use ExecStartPre to configure
# the interface and ExecStopPost to tear it down. BITRATE can be overridden
# via an Environment= variable if desired.
ExecStart=/bin/bash -c 'sleep 5; source /home/neonaero/Documents/odrive-controller/myenv/bin/activate; python3 can_runner.py'
Restart=on-failure
RestartSec=3
Environment=PYTHONUNBUFFERED=1
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/neonaero/.Xauthority
Environment=QT_QPA_PLATFORM=xcb


[Install]
WantedBy=multi-user.target
