# RaspberryPi/systemd/odrive-can.service
[Unit]
Description=ODrive can_runner service (starts the controller and tears down can0 on stop)
After=can0-setup.service
Wants=can0-setup.service

[Service]
Type=simple
User=neonaero
WorkingDirectory=/home/neonaero/Documents/odrive-controller/odrive-controller-odrive-rpi/RaspberryPi
# Use the virtualenv python binary and absolute path to the script
ExecStart=/home/neonaero/Documents/odrive-controller/myenv/bin/python3 /home/neonaero/Documents/odrive-controller/odrive-controller-odrive-rpi/RaspberryPi/can_runner.py
Environment=PYTHONUNBUFFERED=1
Restart=on-failure
# After the service stops, ensure the CAN interface is taken down.
ExecStopPost=/sbin/ip link set can0 down || true
TimeoutStartSec=10

[Install]
WantedBy=multi-user.target
