# RaspberryPi/systemd/odrive-can.service
[Unit]
Description=ODrive can_runner service (starts the controller and tears down can0 on stop)
After=multi-user.target network.target dev-gpiomem.device dev-spidev0.0.device
Requires=network.target

[Service]
Type=simple
User=neonaerolab
WorkingDirectory=/home/neonaerolab/Documents/odrive-controller/RaspberryPi
# Bring up can0 only while this service runs. Use ExecStartPre to configure
# the interface and ExecStopPost to tear it down. BITRATE can be overridden
# via an Environment= variable if desired.
ExecStartPre=sudo /usr/sbin/ip link set can0 type can bitrate 1000000
ExecStartPre=sudo /usr/sbin/ip link set can0 up
ExecStart=/home/neonaerolab/Documents/myenv/bin/python3 can_runner.py
ExecStopPost=sudo /usr/sbin/ip link set can0 down
Restart=on-failure
RestartSec=3
Environment=PYTHONUNBUFFERED=1
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/neonaerolab/.Xauthority
Environment=QT_QPA_PLATFORM=xcb


[Install]
WantedBy=multi-user.target
